# scripts/com.impress.uitests.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.impress.uitests</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tabel/Projects/impress-apps/scripts/start-ui-tests.sh</string>
        <string>all</string>
    </array>

    <!-- Run at 2:00 AM daily -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/uitests-launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/uitests-launchd.log</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
